# backend/docker/supervisord.conf
; Runs the whole backend in one always-on container: web (php-fpm + nginx) plus the two
; long-running Laravel processes that DRIVE all automation. Previously these were a fragile
; backgrounded subshell in start.sh with no restart or visibility, so when it stopped, every
; scheduled job (SEO autopilot, import completion, social, market) silently stopped with it.
; supervisord runs each as a supervised process, restarts it if it dies, and streams logs to
; the container stdout so Render captures them.
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0

; Enable `supervisorctl status` inside the Render shell for at-a-glance process health.
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface.make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; The scheduler: schedule:work is a long-running foreground process that fires every due job
; each minute (SEO 02:00, social 08:30, market 05:30, the every-minute heartbeat + catch-up,
; queue tick). Replaces the old "sleep 60; schedule:run" subshell.
[program:scheduler]
command=php /app/artisan schedule:work
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; The queue worker: processes queued jobs (imported-draft completion, market suggestions,
; etc.) within seconds instead of waiting for a scheduled sweep. --max-time recycles the
; process hourly so long-lived leaks can't accumulate; supervisord restarts it immediately.
[program:queue]
command=php /app/artisan queue:work --sleep=3 --tries=1 --max-time=3600 --backoff=10
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
